# piratemodem-common/files/etc/uci-defaults/11_piratemodem
#!/bin/sh
uci -q batch << EOF
set system.@system[0].hostname='PirateModem'
set network.lan.proto='dhcp'
EOF
